# fs2/usr/lib/ipkg/info/procmon-pc202.list
/etc/rcS.d/S85startprocmon
/opt/ipaccess/bin/procmon
/opt/ipaccess/bin/procmon.491.44.1
/opt/ipaccess/procmon.conf
/etc/init.d/startprocmon
